# quadlets/pathfinder-api.container
[Unit]
Description=Pathfinder API server
After=pathfinder-db.service
After=pathfinder-redis.service
Requires=pathfinder-db.service
Requires=pathfinder-redis.service

[Container]
ContainerName=pathfinder-api
Image=localhost/pathfinder-api:latest
Network=pathfinder.network
# Change the host port (left side) to use a different port
PublishPort=8000:8000

# API keys and app settings loaded from .env
EnvironmentFile=%h/.config/pathfinder/.env

# Infrastructure URLs — must use container network hostnames, override .env
Environment=DATABASE_URL=postgresql+asyncpg://postgres:postgres@pathfinder-db:5432/pathfinder
Environment=REDIS_URL=redis://pathfinder-redis:6379/0
# Uncomment to use a local Ollama instance on the host
# Environment=OLLAMA_BASE_URL=http://host.containers.internal:11434/v1

HealthCmd=curl -fsS http://localhost:8000/health || exit 1
HealthInterval=30s
HealthTimeout=10s
HealthStartPeriod=10s
HealthRetries=3

[Service]
Restart=always

[Install]
WantedBy=default.target
